# termostat-usb-data-uart/files/dserr18usba_init
#!/bin/sh /etc/rc.common
# Copyright (C) 2006 OpenWrt.org

START=99



start() {

	dserr18usba &

}

stop() {
	killall -9 dserr18usba >/dev/null 2>/dev/null
	rm /var/zawis.us &> /dev/null
}
